# .supervisord/celery_worker.conf
[program:bisti-celery-worker]
directory=/code/
process_name=%(program_name)s_%(process_num)02d
command=celery -A apps.core.celery worker --loglevel=info -E
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
user=root
